# vpn-bot.service
[Unit]
Description=VPN Seller Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vpn-bot
ExecStart=/root/vpn-bot/venv/bin/python3 /root/vpn-bot/bot.py
Restart=always
RestartSec=10
Environment=PATH=/root/vpn-bot/venv/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
